# apk/init.d/iknite-config
#!/sbin/openrc-run

name="Iknite Kubernetes configurer"
description="Configure Kubernetes after initialization"


depend() {
    need iknite-init
    after iknite-init
    after kubelet
}

start() {
    ebegin "Running iknite config"
    /sbin/iknite -v debug configure --cluster-check-retries 10 --cluster-check-ok-responses 2 >>/var/log/iknite.log 2>&1
    eend $?
}
