# apk/init.d/iknite-config
#!/sbin/openrc-run

name="Iknite Kubernetes configurer"
description="Configure Kubernetes after initialization"

IKNITE_VERBOSITY="${IKNITE_VERBOSITY:-info}"
IKNITE_CLUSTER_CHECK_RETRIES="${IKNITE_CLUSTER_CHECK_RETRIES:-10}"
IKNITE_CLUSTER_CHECK_OK_RESPONSES="${IKNITE_CLUSTER_CHECK_OK_RESPONSES:-2}"
IKNITE_KUSTOMIZE_DIRECTORY="${IKNITE_KUSTOMIZE_DIRECTORY:-/etc/iknite.d}"

depend() {
    need iknite-init
    after iknite-init kubelet
}

start() {
    ebegin "Running iknite config"
    /usr/bin/env $IKNITE_ENV /sbin/iknite -v $IKNITE_VERBOSITY configure \
        --cluster-check-retries $IKNITE_CLUSTER_CHECK_RETRIES \
        --cluster-check-ok-responses $IKNITE_CLUSTER_CHECK_OK_RESPONSES \
        --kustomize-directory $IKNITE_KUSTOMIZE_DIRECTORY \
        ${IKNITE_OPTS} >>/var/log/iknite.log 2>&1
    eend $?
}
