# apk/init.d/iknite-init
#!/sbin/openrc-run

name="Iknite Kubernetes starter"
description="Ensure the environment is correct and start kubeadm"

IKNITE_VERBOSITY="${IKNITE_VERBOSITY:-info}"

depend() {
    after net
    need containerd
    # before kubelet
}

start() {
    ebegin "Running iknite start"

    $IKNITE_ENV /sbin/iknite -v $IKNITE_VERBOSITY start \
        --openrc \
        ${IKNITE_OPTS} \
        >>/var/log/iknite.log 2>&1
    eend $?
}
