# apk/init.d/iknite-init
#!/sbin/openrc-run

name="Iknite Kubernetes starter"
description="Ensure the environment is correct and start kubeadm"

depend() {
    after net
    need crio
    # before kubelet
}

start() {
    ebegin "Running iknite start"
    /sbin/iknite -v debug start --openrc >>/var/log/iknite.log 2>&1
    eend $?
}
