# packaging/apk/incus-agent/init.d-incus-agent-setup
#!/sbin/openrc-run

description="Incus - agent - setup"
command=/usr/bin/incus-agent-setup
required_dirs=/dev/virtio-ports/
